# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:vllm]
command=python3 run.py
directory=/vllm-service
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
startsecs=10
startretries=3

[program:chatbot]
command=python3 -m uvicorn app.main:app --host 0.0.0.0 --port 8300
directory=/chatbot-service
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
startsecs=5
startretries=3

[program:multimodal]
command=python3 -m uvicorn app.main:app --host 0.0.0.0 --port 8500
directory=/multimodal-service
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
startsecs=5
startretries=3
